# dev/supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=cron -f

[program:php]
command=php-fpm -R
